# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/acl.list
/usr/sbin/aclctl
/usr/sbin/aclhijackdns
/usr/sbin/acl_update_name
/usr/sbin/acld
/etc/init.d/acl
