# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/bin/login
/lib/cfgmgr/cfgmgr.sh
/etc/init.d/net-wan
/etc/resolv.conf
/etc/init.d/done
/var
/sbin/cfgapply
/lib/cfgmgr/opmode.sh
/etc/hotplug2-init.rules
/sbin/check_status.sh
/etc/preinit
/etc/config/tt3
/sbin/mount_root
/lib/upgrade/platform.sh
/sbin/cmd_sipalg
/etc/hotplug2-common.rules
/etc/shells
/etc/protocols
/etc/netgear.cfg
/sbin/detect_factory
/etc/language.cfg
/etc/hotplug.d/button/00-button
/etc/dnidata.cfg
/sbin/cfgmgr/cfg_email_ntp
/etc/init.d/init.enet
/etc/init.d/powerctl
/etc/uci-defaults/min_free_kbytes
/etc/fstab
/sbin/printhosts
/sbin/led.sh
/etc/banner
/etc/init.d/net-lan
/etc/config/nowan
/etc/rc.common
/init
/etc/uci-defaults/led
/sbin/check_board_parameter
/etc/init.d/watchdog
/sbin/power_led.sh
/lib/functions/uci-defaults.sh
/bin/fbwifi
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/init.d/usb
/lib/ufsd/ufsd.ko
/sbin/cfgmgr/cfg_apmode
/sbin/usb_disk_event
/etc/services
/sbin/cloud
/etc/sysupgrade.conf
/lib/preinit/03_preinit_do_ipq806x.sh
/sbin/cmd_autoreset
/etc/init.d/bridge
/lib/upgrade/keep.d/base-files
/etc/uci-defaults/migrate-sysctl
/sbin/wan-ipup
/etc/functions.sh
/sbin/firstboot
/etc/config/lanwan
/sbin/cfgmgr/cfg_opmode
/etc/hosts
/etc/uci-defaults/migrate-shadow
/bin/ookla
/etc/hotplug.d/usb/10-usb
/lib/functions.sh
/lib/ipq806x.sh
/sbin/cfgmgr/cfg_extmode
/etc/init.d/rcS
/sbin/wps_pbc
/rom/note
/etc/inittab
/sbin/qos.sh
/etc/openwrt_release
/etc/uci-defaults/syslog
/sbin/wifi
/etc/config/swinit
/sbin/run-ramfs
/sbin/cmd_ddns
/etc/profile
/sbin/cpuutil
/etc/uci-defaults/buttons
/etc/TZ
/lib/functions/service.sh
/etc/init.d/umount
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/mnt
/etc/mtab
/etc/init.d/boot
/lib/functions/gpio-button-hotplug.sh
/sbin/cfgmgr/cfg_brmode
/lib/functions/boot.sh
/bin/ipcalc.sh
/sbin/kwilt
/etc/diag.sh
/etc/uci-defaults/network
