# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/bridge-mode.list
/etc/init.d/br-mode
/usr/sbin/br-dhcp-helper
/usr/sbin/start_br_udhcpc
