# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/detcable.list
/etc/init.d/detcable
/usr/bin/detcable
/usr/bin/flushRoute
