# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/forked-daapd.list
/usr/lib/forked-daapd/forked-daapd-sqlext.la
/usr/sbin/forked-daapd
/lib/upgrade/keep.d/forked-daapd
/usr/lib/forked-daapd/forked-daapd-sqlext.so
/etc/forked-daapd.conf
/etc/init.d/forked-daapd
/usr/sbin/itunes_allow_control
