# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/guest_dhcp.list
/etc/hotplug.d/iface/35-guest_dhcp
/etc/dhcp.guest.conf
/etc/init.d/guest_dhcp
