# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmpproxy
/sbin/cmdigmp
/etc/config/igmpproxy
/etc/hotplug.d/iface/50-igmpproxy
/sbin/igmpproxy
