# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/kmod-dni-mtdoops.list
/etc/init.d/dni-debug.init
/sbin/wireless_log.sh
/sbin/console_log.sh
/sbin/mirror.sh
/sbin/capture_packet.sh
/lib/modules/3.4.103/mtdoops.ko
/sbin/debug_save_panic_log
/sbin/basic_log.sh
/sbin/enable_collect_console_log.sh
/sbin/save_collect_console_log.sh
/sbin/wlandebug.sh
/sbin/debug_dutinfo.sh
/sbin/debug_telnetenable.sh
