# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/kmod-qca-nss-ecm.list
/etc/sysctl.d/qca-nss-ecm.conf
/etc/uci-defaults/99-qca-nss-ecm
/usr/bin/ecm_dump.sh
/etc/init.d/qca-nss-ecm
/etc/firewall.d/qca-nss-ecm
/lib/modules/3.4.103/ecm.ko
