# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/kmod-qca-wifi.list
/lib/modules/3.4.103/umac.ko
/lib/modules/3.4.103/ath_hal.ko
/etc/config/ssid-steering.config
/usr/sbin/thermaltool
/usr/sbin/athadhoc
/lib/modules/3.4.103/ath_rate_atheros.ko
/usr/sbin/apstats
/usr/sbin/spectraltool
/lib/modules/3.4.103/smart_antenna.ko
/usr/sbin/assocdenialnotify
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/lib/modules/3.4.103/asf.ko
/lib/modules/3.4.103/wlan_ext.ko
/usr/sbin/radartool
/etc/init.d/madwifi
/lib/wifi/qcawifi_countrycode.txt
/lib/modules/3.4.103/adf.ko
/usr/sbin/wifitool
/etc/init.d/ssid_steering.init
/usr/lib/libtlvparser.so
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/usr/sbin/tx99tool
/lib/modules/3.4.103/hst_tx99.ko
/etc/config/wireless
/lib/wifi/33-qca-wifi
/usr/sbin/athstats
/usr/lib/libtlvtemplate.so
/lib/modules/3.4.103/ath_dfs.ko
/usr/lib/libtlvcmdrsp.so
/sbin/wifi_try
/lib/sbin/wifi_try
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/modules/3.4.103/ath_pktlog.ko
/lib/modules/3.4.103/ath_spectral.ko
/usr/sbin/athdiag
/lib/modules/3.4.103/ath_dev.ko
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/usr/sbin/wps_enhc
/usr/sbin/pktlogconf
