# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/kmod-kcode.list
/lib/modules/3.4.103/GPL_NetUSB.ko
/lib/modules/3.4.103/NetUSB.ko
/sbin/KC_PRINT
/sbin/KC_BONJOUR
/etc/init.d/kcode
