# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/openvpn-openssl.list
/etc/openvpn/push_routing_rule
/etc/init.d/zzprefix-check_cert_files
/etc/init.d/openvpn
/usr/sbin/openvpn
/etc/rc.d/S60openvpn
/lib/upgrade/keep.d/openvpn-openssl
/etc/openvpn/download
/lib/upgrade/keep.d/openvpn
