# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/pot.list
/etc/init.d/pot
/usr/sbin/potd
/usr/sbin/ntpst
/usr/sbin/stamac
/usr/sbin/potval
