# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/streamboost-misc.list
/usr/local/bin/jiggle_firewall
/etc/config/system
/etc/sysctl.conf
/etc/ssl/certs/d94217ed.0
/usr/local/bin/reset_wan
/etc/config/firewall
/etc/rc.local
/usr/local/bin/stfu
/etc/ssl/certs/d77c6e41.0
/etc/crontabs/root
/etc/config/fstab
/etc/init.d/sysfs_config
/usr/local/bin/apply_appflow
/etc/init.d/uci2streamboost
