# build_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/root-ipq806x/usr/lib/opkg/info/wlan-common.list
/etc/hotplug.d/wps/00-wps
/sbin/update-wifi
/etc/hotplug.d/wlanlog/00-wlanlog
/sbin/wlandebug
/lib/wifi/dniconfig.sh
/etc/hotplug.d/wps/01-wps-led
/etc/wifi_config_comp.awk
/etc/dni-wifi-config
/sbin/wpsled
/etc/init.d/wlan-common
/etc/ath/wifi.conf
/sbin/wlan
