# staging_dir/target-arm_v7-a_uClibc-0.9.33.2_eabi/packages/shortcut-fe.list
./
./host
./etc
./etc/init.d
./etc/init.d/shortcut-fe
./usr
./usr/bin
./usr/bin/sfe_dump
./usr/include
./usr/include/shortcut-fe
./usr/include/shortcut-fe/sfe.h
